# dataserver/sv/zotero-download/run
#!/bin/sh

cd /srv/zotero/dataserver/processor/download
exec 2>&1
exec chpst -u www-data:www-data php5 daemon.php
